# deploy/systemd/visp-server.service
[Unit]
Description=VISP API
After=network-online.target postgresql.service tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/visp/apps/server
EnvironmentFile=/etc/visp/app.env
ExecStart=/usr/local/bin/bun run dist/index.mjs
Restart=always
RestartSec=2
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
